# deploy/hibs-trading-daemon.service
[Unit]
Description=HIBS racing trading daemon (simulation sandbox)
After=network.target hibs-racing.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/hibs-racing
Environment=HOME=/opt/hibs-racing
Environment=PYTHONPATH=/opt/hibs-racing/src
Environment=HIBS_RACING_DB_PATH=/mnt/hibs-ramdisk/feature_store.sqlite
Environment=HIBS_LIVE_TRADING_ENABLED=false
Environment=HIBS_LIQUIDITY_ROUTER_ACTIVE=false
EnvironmentFile=-/opt/hibs-racing/.env
ExecStart=/opt/hibs-racing/.venv/bin/hibs-racing trading-daemon
Restart=on-failure
RestartSec=15
NoNewPrivileges=true
PrivateTmp=true
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
